# deploy/ask-territory-refresh.timer
# Runs ask-territory-refresh.service every 6 hours.
[Unit]
Description=Run Ask Territory refresh every 6 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
